# modules/tomoconnor-postgres/files/check_lag.cron
# /etc/cron.d/check_lag
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

* * * * *	root	/usr/local/bin/check_lag.sh

#EOF
