# modules/tomoconnor-graylog/files/graylog2-web.conf
#upstart for graylog2-web
start on started graylog2
stop on [06]
exec /opt/oss/graylog2/0.9/web/script/rails server -e production

